# projects/managers/check_critical.path
[Unit]
Description=Watch My File

[Path]
PathModified=/path/to/watched/file
Unit=check_critical.service

[Install]
WantedBy=multi-user.target
