# projects/managers/check_critical.service
[Unit]
Description=Check File for Critical Error

[Service]
Type=simple
ExecStart=/path/to/your/check_script.sh
