# projects/managers/pass_monitor.path
[Unit]
Description="Monitor the /tmp/passwd file for changes"

[Path]
PathModified=/tmp/passwd
Unit=pass_monitor.service

[Install]
WantedBy=multi-user.target
